# docker/crontab
# Daily snapshot cron job - runs at 00:05 UTC
5 0 * * * cd /app && python scripts/daily_snapshot.py --yesterday >> /app/logs/daily_snapshot.log 2>&1
# Empty line required at end of crontab